# proxy-server/deploy/obsidian-proxy.service
[Unit]
Description=Obsidian Claude Proxy Server
After=network.target

[Service]
Type=simple
User=fxcyf
WorkingDirectory=/home/fxcyf/MyFiles/random/obsidian-ai-daily/proxy-server
ExecStart=/home/fxcyf/.nvm/versions/node/v20.19.5/bin/node dist/server.js
Restart=always
RestartSec=5
Environment=PORT=27090
Environment=AUTH_TOKEN=986753421
Environment=MCP_CONFIG=/home/fxcyf/MyFiles/random/obsidian-ai-daily/proxy-server/mcp-config.json
Environment=TOOL_POLICY_PATH=/home/fxcyf/MyFiles/random/obsidian-ai-daily/agent-tool-policy.json
Environment=CLAUDE_MODEL=sonnet
Environment=CLAUDE_PATH=/home/fxcyf/.local/bin/claude
Environment="VAULT_PATH=/home/fxcyf/MyFiles/Obsidian/All in One"
Environment=PATH=/home/fxcyf/.local/bin:/home/fxcyf/.nvm/versions/node/v20.19.5/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
